# deploy/systemd/kitchorify-backup-sync.timer
[Unit]
Description=Run Kitchorify off-host backup sync daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
